# docs/zaomeng-web.service
[Unit]
Description=Zaomeng Next.js Web App
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Downloads/zaomeng/project/projects
Environment=NODE_ENV=production
Environment=PORT=5000
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/Downloads/zaomeng/project/projects/.coze-logs/systemd-web.log
StandardError=append:/home/ubuntu/Downloads/zaomeng/project/projects/.coze-logs/systemd-web-error.log

[Install]
WantedBy=multi-user.target
